# docker/backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ---------------------------------------------------------------------------
# ocr-api  (OpenAI Vision OCR endpoint — must start before ocr-service)
# Other services that depend on it start with slight delays via startdelay.
# ---------------------------------------------------------------------------
[program:ocr-api]
command=uvicorn main:app --host 0.0.0.0 --port 8002
directory=/app/ocr_api
environment=OTEL_SERVICE_NAME="ocr-api"
stdout_logfile=/var/log/supervisor/ocr-api.log
stderr_logfile=/var/log/supervisor/ocr-api.log
autostart=true
autorestart=true
priority=10

# ---------------------------------------------------------------------------
# ingest-api  (FastAPI, port 8000)
# ---------------------------------------------------------------------------
[program:ingest-api]
command=/bin/sh -c ". /app/ingest_api/set-env.sh && exec uvicorn app.main:app --host 0.0.0.0 --port 8000"
directory=/app/ingest_api
environment=REDIS_URL="%(ENV_INGEST_REDIS_URL)s",OTEL_SERVICE_NAME="ingest-api"
stdout_logfile=/var/log/supervisor/ingest-api.log
stderr_logfile=/var/log/supervisor/ingest-api.log
autostart=true
autorestart=true
priority=20

# ---------------------------------------------------------------------------
# retrieval-api  (FastAPI, port 8001)
# ---------------------------------------------------------------------------
[program:retrieval-api]
command=/bin/sh -c ". /app/retrieval_api/set-env.sh && exec uvicorn app.main:app --host 0.0.0.0 --port 8001"
directory=/app/retrieval_api
environment=REDIS_URL="%(ENV_RETRIEVAL_REDIS_URL)s",OTEL_SERVICE_NAME="retrieval-api"
stdout_logfile=/var/log/supervisor/retrieval-api.log
stderr_logfile=/var/log/supervisor/retrieval-api.log
autostart=true
autorestart=true
priority=20

# ---------------------------------------------------------------------------
# ingestion-worker
# ---------------------------------------------------------------------------
[program:ingestion-worker]
command=/bin/sh -c ". /app/ingestion_worker/set-env.sh && exec python -m app.main"
directory=/app/ingestion_worker
environment=REDIS_URL="%(ENV_INGEST_REDIS_URL)s",OTEL_SERVICE_NAME="ingestion-worker"
stdout_logfile=/var/log/supervisor/ingestion-worker.log
stderr_logfile=/var/log/supervisor/ingestion-worker.log
autostart=true
autorestart=true
priority=20

# ---------------------------------------------------------------------------
# ocr-service  (RabbitMQ consumer — calls ocr-api via localhost:8002)
# ---------------------------------------------------------------------------
[program:ocr-service]
command=/bin/sh -c ". /app/ocr_service/set-env.sh && exec python -m app.main"
directory=/app/ocr_service
environment=OCR_API_URL="http://localhost:8002/ocr",OTEL_SERVICE_NAME="ocr-service"
stdout_logfile=/var/log/supervisor/ocr-service.log
stderr_logfile=/var/log/supervisor/ocr-service.log
autostart=true
autorestart=true
priority=20

# ---------------------------------------------------------------------------
# embedding-service
# ---------------------------------------------------------------------------
[program:embedding-service]
command=/bin/sh -c ". /app/embedding_service/set-env.sh && exec python -m app.main"
directory=/app/embedding_service
environment=REDIS_URL="%(ENV_EMBEDDING_REDIS_URL)s",OTEL_SERVICE_NAME="embedding-service"
stdout_logfile=/var/log/supervisor/embedding-service.log
stderr_logfile=/var/log/supervisor/embedding-service.log
autostart=true
autorestart=true
priority=20
